# instance/services/jupyter/run
#!/bin/sh

exec su jupyter -s /bin/sh -c "cd /home/jupyter/; jupyter notebook --no-browser"
